# music_assistant_proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Get config values
export SERVER_HOST=$(bashio::config 'server_host')
export SERVER_PORT=$(bashio::config 'server_port')
export MA_TOKEN=$(bashio::config 'ma_token')

bashio::log.info "Starting Music Assistant Proxy..."
bashio::log.info "Server Host: ${SERVER_HOST}"
bashio::log.info "Server Port: ${SERVER_PORT}"
bashio::log.info "Music Assistant Token: ${MA_TOKEN}"

# Generate nginx config from template
envsubst '$SERVER_HOST $SERVER_PORT $MA_TOKEN' < /etc/nginx/nginx.conf > /tmp/nginx.conf
mv /tmp/nginx.conf /etc/nginx/nginx.conf

# Create required directories
mkdir -p /var/log/nginx
mkdir -p /var/lib/nginx/tmp
chown nginx:nginx /var/lib/nginx/tmp

bashio::log.info "Starting Nginx..."
# Start Nginx
exec 2>&1
exec nginx -g 'daemon off;'
